# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# --- XRDP (RDP access on 3389) — run in foreground ---
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp.log
stderr_logfile=/var/log/xrdp.err

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp-sesman.log
stderr_logfile=/var/log/xrdp-sesman.err

# --- Xvnc (XFCE session on :1) in foreground for noVNC ---
[program:vnc]
command=Xvnc :1 -geometry 1280x800 -depth 24 -SecurityTypes None -rfbport 5901 -pn -AlwaysShared -AcceptKeyEvents -AcceptPointerEvents
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/vnc.log
stderr_logfile=/var/log/vnc.err

[program:xfce4]
command=bash -lc "export DISPLAY=:1; exec startxfce4"
autostart=true
autorestart=true
priority=25
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err

# --- noVNC web bridge on Railway $PORT ---
[program:novnc]
command=bash -lc "websockify --web=/usr/share/novnc %(ENV_PORT)s localhost:5901"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
