# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# --- XRDP (RDP on 3389) — fork=false in xrdp.ini keeps it foreground ---
[program:xrdp]
command=/usr/sbin/xrdp
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp.log
stderr_logfile=/var/log/xrdp.err

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp-sesman.log
stderr_logfile=/var/log/xrdp-sesman.err

# --- Xvnc (TigerVNC) in foreground on display :1, port 5901 ---
[program:xvnc]
command=Xtigervnc :1 -geometry 1280x800 -depth 24 -SecurityTypes None -rfbport 5901 -pn -fg -AlwaysShared
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/xvnc.log
stderr_logfile=/var/log/xvnc.err

# --- XFCE4 desktop on display :1 ---
[program:xfce4]
command=bash -lc "export DISPLAY=:1; exec startxfce4"
autostart=true
autorestart=true
priority=25
environment=DISPLAY=":1",XDG_SESSION_TYPE="x11"
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err

# --- noVNC web bridge (browser access) ---
[program:novnc]
command=bash -lc "websockify --web=/usr/share/novnc %(ENV_PORT)s localhost:5901"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
