# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# --- XRDP (RDP access on 3389) ---
[program:xrdp]
command=/usr/sbin/xrdp
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp.log
stderr_logfile=/var/log/xrdp.err

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xrdp-sesman.log
stderr_logfile=/var/log/xrdp-sesman.err

# --- VNC server (XFCE session on :1) for noVNC ---
[program:vnc]
command=bash -lc "vncserver :1 -geometry 1280x800 -depth 24 -SecurityTypes None"
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/vnc.log
stderr_logfile=/var/log/vnc.err

# --- noVNC web bridge on Railway $PORT ---
[program:novnc]
command=bash -lc "websockify --web=/usr/share/novnc %(ENV_PORT)s localhost:5901"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

# --- Keep the session environment tidy (optional) ---
[program:cleanup]
command=bash -lc "sleep 30 && echo booted"
autostart=true
autorestart=false
